# docker/rootfs/etc/services.d/nodered/run
#!/command/with-contenv sh
set -eu
exec node-red --userDir /data/nodered
