# docker/rootfs/etc/services.d/openplc/run
#!/command/with-contenv sh
set -eu
exec sh -c 'echo "OPENPLC_RUNTIME_PLACEHOLDER: implement start command" && sleep infinity'
