# docker/rootfs/etc/services.d/cron/run
#!/command/with-contenv sh
set -eu
exec cron -f
