# docker/rootfs/etc/services.d/cron/run
#!/command/with-contenv sh
set -eu
exec crond -f -l 8
